# etc/smartmeter.service.default
[Unit]
Description=SmartMeter
Documentation=https://gitlab.com/smart-home-tools/smartmetertomqtt

[Service]
Type=simple
ExecStart=/usr/local/bin/SmartMeterToMqtt -f /etc/smartmeter.json
#User=pi
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
