# etc/smartmeter.service.default
[Unit]
Description=SmartMeter
Documentation=https://github.com/oetken/smartmetertomqtt

[Service]
Type=simple
ExecStart=SmartMeterToMqtt -f /etc/smartmeter.json
#User=pi
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
